# supervisord.conf
# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:flask]
command=flask run --host=0.0.0.0 --port=5000
directory=/app
autorestart=true
environment=FLASK_APP="app.py"  # Replace with your Flask entry file (e.g., app.py)